# puppet/development/modules/mailhog/files/mailhog.conf
# mailhog - Web and API based SMTP testing

description     "MailHog server"

start on startup
task
exec /bin/sh /opt/mailhog/mailhog_start.sh